# crontab
* * * * * python /usr/src/selenium_auto_test/test_website.py >> /var/log/cron/cron.log 2>&1
